# www/python/racetiming.conf
start on runlevel [2345]
stop on runlevel [016]

respawn
exec python /var/www/html/python/racetiming10_dev.py
